# conf/supervisor.conf
# -*- conf -*-

[program:detroitnow]
command=/srv/sites/detroitnow/conf/startup
autostart=true
autorestart=true
stderr_logfile=/srv/log/detroitnow/supervisor_err.log
numprocs=1
